# devices/router/openwrt/files/etc/rc.local
# Setup CPU frequency scaling
echo ondemand > /sys/devices/system/cpu/cpufreq/policy0/scaling_governor
