# devices/router/openwrt/files/etc/hotplug.d/iface/10-sfp
#!/bin/sh

[ "$ACTION" = ifup -a "$DEVICE" = eth2 ] || exit 0

logger -t sfp "Force 1Gbps full-duplex on SFP module (eth2)"
ethtool -s eth2 speed 1000 duplex full autoneg off
