# devices/router/openwrt/files/etc/hotplug.d/iface/10-sfp
#!/bin/sh

[ "$ACTION" = ifup -a "$DEVICE" = sfp-wan ] || exit 0

logger -t sfp "Force 1Gbps full-duplex on SFP module (sfp-wan)"
ethtool -s sfp-wan speed 1000 duplex full autoneg off
